# system_files/desktop/shared/usr/lib/systemd/user/bazzite-user-setup.service
[Unit]
Description=Configure zeropoint for current user

[Service]
Type=simple
ExecStart=/usr/libexec/zeropoint-user-setup

[Install]
WantedBy=default.target
